# roles/branding/templates/etc/update-motd.d/20-logo.j2
#!/bin/sh
echo
cat /etc/{{ ianseobox_logo_file }}
